# system/doppel2_webserver.ini
[program:doppel2_webserver]
command=/usr/bin/python2 /srv/doppel2/manage.py runserver 127.0.0.1:8000
user=http
umask=022
redirect_stderr=true
stopasgroup=true
