# system/doppel2_webserver.ini
[program:chain_webserver]
command=/usr/bin/python2 /srv/chain-api/manage.py runserver 127.0.0.1:8000
user=http
umask=022
redirect_stderr=true
stopasgroup=true
