# system/etc/supervisor/conf.d/chainsocketpolicyd.conf
[program:chainsocketpolicyd]
command=chainsocketpolicyd
user=root
redirect_stderr=true
