# system/etc/supervisor/conf.d/medialab_rfidpost.conf
[program:medialab_rfidpost]
command=rfidpost http://chain-api.media.mit.edu/sites/5
user=www-data
umask=022
redirect_stderr=true
stopasgroup=true
