# system/etc/supervisor/conf.d/tidmarsh_tidpost.conf
[program:tidmarsh_tidpost]
command=tidpost http://chain-api.media.mit.edu/sites/2
user=www-data
umask=022
redirect_stderr=true
stopasgroup=true
